# target/linux/ipq40xx/base-files/etc/init.d/bootcomplete
#!/bin/sh /etc/rc.common

START=99

boot() {
	case $(board_name) in
	sony,ncp-hg100)
		# boot complete LED blink off
		echo 0xb0 > /sys/class/i2c-dev/i2c-1/device/1-0015/led
		# front touch panel enable
		i2cset -y 1 0x14 0x0f 1
		;;
	esac
}
