# target/linux/kirkwood/nas/base-files/etc/init.d/done_led
#!/bin/sh /etc/rc.common

START=96
boot() {
	case "$(board_name)" in
	iodata,hdl-a-sata|\
	iodata,hdl-a-usb|\
	iodata,hdl2-a-sata|\
	iodata,hdl2-a-usb)
		# green, on
		echo ":sts on" > /dev/ttyS1
		;;
	esac
}
